# phs-api.service
[Unit]
Description=PHS API Server
After=network.target postgresql.service

[Service]
Type=simple
User=mapdroid
WorkingDirectory=/home/mapdroid/phs-api
ExecStart=/usr/bin/node api-server.js

# Auto-restart on crash
Restart=always
RestartSec=5s

# Graceful shutdown (15 sec timeout)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

# Resource limits
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phs-api

# Load environment from file (deploy to /etc/default/phs-api)
EnvironmentFile=/etc/default/phs-api

[Install]
WantedBy=multi-user.target
